# supervisord.conf
[supervisord]
nodaemon=true

[program:duckling]
directory=/duckling
command=./duckling-example-exe

[program:flask]
directory=/app
command=python3 app.py